# src/lib/systemd/user/flatalias.service
[Unit]
Description=Run Flatalias to update cached command line aliases of Flatpaks
After=default.target

[Service]
Type=oneshot
ExecStart=@bindir@/flatalias %S/flatalias.profile

[Install]
WantedBy=default.target
