# src/lib/systemd/user/flatalias.path
[Unit]
Description=Watch Flatpak exports and run Flatalias to update cached command line aliases of Flatpaks

[Path]
PathChanged=/var/lib/flatpak/exports/bin
PathChanged=%D/flatpak/exports/bin

[Install]
WantedBy=default.target
